# supervisord.conf
[program:php-fpm]
  command=/usr/sbin/php-fpm7.4 -c /etc/php/7.4/fpm/php-fpm.conf
  priority=1
  startretries=10

[program:nginx]
  command=/usr/sbin/nginx
  priority=2
  startretries=10
